# ansible/roles/prometheus-node/templates/init.conf.j2
start on startup

script
  /opt/prometheus-node/node_exporter
end script
