# ansible/roles/prometheus/templates/init.conf.j2
start on startup

script
  /usr/bin/prometheus --config.file=/opt/prometheus/prometheus.yml
end script
